# backend/.settings/supervisord_server.conf
[supervisord]
user=root
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:jheep-server]
user=root
command=jheep server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:dramatiq-worker]
user=root
command=jheep worker -p 4 -t 2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
